# package/mtk/applications/luci-app-wifi-mtk/root/etc/uci-defaults/luci-app-wifi-mtk.default
#!/bin/sh

if [ -d /www/luci-static/resources/view/network ]; then
	mv -f /usr/share/luci-app-wifi-mtk/wireless-mtk.js /www/luci-static/resources/view/network/wireless.js
fi

# if [ -d /usr/share/luci/menu.d ]; then
# 	mv -f /usr/share/luci-app-wifi-mtk/luci-mod-status.json /usr/share/luci/menu.d/luci-mod-status.json
# fi
#
if [ -f /usr/share/luci/menu.d/luci-app-nlbwmon.json ]; then
	mv -f /usr/share/luci-app-wifi-mtk/luci-app-nlbwmon.json /usr/share/luci/menu.d/luci-app-nlbwmon.json
fi

if [ -f /etc/wireless/l1profile.dat ]; then
	sed -i 's|wifi/mtwifi.lua|netifd/wireless/mtk.sh|g' /etc/wireless/l1profile.dat
	sed -i 's|mtwifi;mtwifi|mtk;mtk|g' /etc/wireless/l1profile.dat
fi

luci_network_js="/www/luci-static/resources/network.js"

if [ -f $luci_network_js ]; then
	if ! grep -sq 'ra\[xiyez\]' $luci_network_js ; then
		sed -i 's#/^wlan\\d+/,#/^wlan\\d+/,/^ra[xiyez]?\\d+/,/^apcli[xiyez]?\\d+/,/^wds[xiyez]?\\d+/,/^mesh[xiyez]?\\d+/,#' $luci_network_js
	fi
fi
